# deployment/ohcampus-backend.service
[Unit]
Description=OhCampus Counselor Backend API
After=network.target mongodb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/counselor.ohcampus.com/backend
Environment="PATH=/var/www/counselor.ohcampus.com/backend/venv/bin"
EnvironmentFile=/var/www/counselor.ohcampus.com/backend/.env
ExecStart=/var/www/counselor.ohcampus.com/backend/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8001 --workers 4
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
